# deploy/systemd/job-ticket-production-backup.timer
[Unit]
Description=Run Job Ticket System production backup daily

[Timer]
OnCalendar=*-*-* 02:15:00
Persistent=true
RandomizedDelaySec=5m
Unit=job-ticket-production-backup.service

[Install]
WantedBy=timers.target
